# docker/etc/supervisor.d/bootstrap.conf
[program:bootstrap]
; checks that container is ready for masabot use and then launches
; masabot

command=/app/kill-on-failure.sh /app/bootstrap.sh
priority=1
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true