# docker/etc/supervisor/programs/bootstrap.conf
[program:bootstrap]
; checks that container is ready for masabot use and then launches
; masabot

command=/apps/kill-on-failure.sh /apps/bootstrap.sh
priority=1
autostart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true